# crontab
*/10 * * * * ~/myskichalet/denormalize.sh
