# crontab
*/10 * * * * ~/myskichalet/denormalize.sh
0 3 * * * ~/myskichalet/interhome.sh
